# scripts/anime_rankings.timer
[Unit]
Description=Run anime_rankings daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
